# src/netboot-base/etc/inittab
# Simple inittab
# $Header: /var/cvsroot/gentoo/src/netboot-base/etc/inittab,v 1.2 2004/10/06 14:50:23 vapier Exp $

# System init
::sysinit:/sbin/rc sysinit

# Main shell
::respawn:-/bin/sh

# What to do at the "Three Finger Salute"
::ctrlaltdel:/sbin/reboot

# System reboot
::shutdown:/sbin/rc reboot
