# src/netboot-base/etc/inittab
# Simple inittab

# System init
::sysinit:/sbin/rc sysinit

# Main shell
::respawn:-/bin/sh

# What to do at the "Three Finger Salute"
::ctrlaltdel:/sbin/reboot

# System reboot
::shutdown:/sbin/rc reboot
