# src/livecd/profiles/rtcw-et/aux-files/md
#!/sbin/runscript
# Copyright 1999-2002 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: /var/cvsroot/gentoo/src/livecd/profiles/rtcw-et/aux-files/Attic/md,v 1.1 2003/10/19 03:46:50 livewire Exp $

HILITE="\033[1m"
NORMAL="\033[0m"

RAID="no"
CMDLINE="`cat /proc/cmdline`"
for x in $CMDLINE
do
        if [ "$x" = "noraid" ]
        then
                RAID="no"
        fi
done

depend() {
	after checkfs
	after *
}

blurb() {
	echo -ne ${HILITE}${1}
	#blurblen=${#1}
}

backup() {
	echo -ne "\033[0G\033[0K"
}
start() {
        ebegin "Raid Setup"

if [ "$RAID" = "yes" ]
        then
                ebegin "Auto-detecting EVMS modules"
                # Disable kernel messages to console:
                echo "0" > /proc/sys/kernel/printk
                local x
                for x in /lib/modules/*/kernel/drivers/evms/*
                do
                        y=${x##*/}
                        y=${y%%.*}
                        echo -ne " ${GOOD}*${NORMAL} Scanning for ${y}..."
                        insmod -f ${x} > /dev/null 2>&1
                        backup
                done
        fi

}

stop() {
	return
}

# vim:ts=4
