# users/pvdabeel/livecd/files/loop/etc/init.d/checkfs
#!/sbin/runscript

# $Header: /var/cvsroot/gentoo/users/pvdabeel/livecd/files/loop/etc/init.d/checkfs,v 1.1.1.1 2004/02/13 00:53:19 pvdabeel Exp $

depend() {
	need checkroot
}

start() {

    ## No need to start repairing a read-only filesystem

    eend 0
    return 0

}


# vim:ts=4
