# src/livecd/profiles/lw-ut2003/aux-files/hwsetup
#!/sbin/runscript
# Copyright 1999-2002 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: /var/cvsroot/gentoo/src/livecd/profiles/lw-ut2003/aux-files/Attic/hwsetup,v 1.1 2003/10/19 03:46:49 livewire Exp $

HILITE="\033[1m"
NORMAL="\033[0m"

DETECT="yes"
CMDLINE="`cat /proc/cmdline`"
for x in $CMDLINE
do
        if [ "$x" = "nodetect" ]
        then
                DETECT="no"
        fi
done

depend() {
	after checkfs
	after local
}

blurb() {
	echo -ne ${HILITE}${1}
	#blurblen=${#1}
}

backup() {
	echo -ne "\033[0G\033[0K"
}
start() {
        ebegin "Hardware detection"

if [ "$DETECT" = "yes" ]
        then
        # Disable kernel messages to console:
        echo "0" > /proc/sys/kernel/printk
		/etc/init.d/autoconfig
		eend
        fi

}

stop() {
	return
}

# vim:ts=4
