# src/patchsets/openafs/scripts/openafs-server
#!/sbin/runscript
# Copyright 1999-2005 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

depend() {
	need net
	use logger dns
}

start() {
	ebegin "Starting OpenAFS BOS server"
	failed=1

	eindent

	if ! /sbin/start-stop-daemon --quiet --start --background --make-pidfile --pidfile /var/run/bosserver.pid --exec /usr/sbin/bosserver -- -nofork ${BOSSERVER_OPTIONS}; then
		eerror "BosServer was already running, continuing without restart"
		eend 1
	else
		eend 0
		failed=0
	fi

	eoutdent

	if [ $failed != 0 ]; then
		return 1
	fi
}

stop() {
	ebegin "Stopping OpenAFS BOS server"
	failed=1

	eindent

	# This kindly kills all server processes
	if ! /usr/bin/bos shutdown localhost -localauth -wait; then
		eend 1
	else
		start-stop-daemon --quiet --stop --pidfile /var/run/bosserver.pid
		failed=$?
		eend $failed
	fi

	eoutdent

	if [ $failed != 0 ]; then
		return 1
	fi
}

