# src/embedded/baselayout-lite/image/etc/init.d/localmount
#!/sbin/runscript
# Copyright 1999-2005 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo/src/embedded/baselayout-lite/image/etc/init.d/Attic/localmount,v 1.2 2005/03/06 06:28:08 iggy Exp $

depend() {
	need checkfs
}

start() {
	local usbfs
	# Mount local filesystems in /etc/fstab.
	ebegin "Mounting local filesystems"
	mount -n -at noproc,noshm,${NO_NET_FS_LIST} >/dev/null
	eend $? "Some local filesystem failed to mount"

	# Make sure we insert usbcore if its a module
	if [ -f /proc/modules ] ; then
		# >/dev/null to hide errors from non-USB users
		modprobe usbcore &> /dev/null
	fi
	
	# Check what USB fs the kernel support.  Currently
	# 2.5+ kernels, and later 2.4 kernels have 'usbfs',
	# while older kernels have 'usbdevfs'.
	if [ "$(grep -q usbfs /proc/filesystems)" -eq "0" ] ; then
		usbfs=usbfs
	elif [ "$(grep -q usbdevfs /proc/filesystems)" -eq "0" ] ; then
		usbfs=usbdevfs
	fi

	if [ -n ${usbfs} ] && \
	   [ -e /proc/bus/usb ] && [ ! -e /proc/bus/usb/devices ]
	then
		ebegin "Mounting USB device filesystem (${usbfs})"
#		# Fetch usb gid from /etc/group; fixes bug 35860
#		usbgid=$(awk -F: '/^usb:/{print $3; exit}' /etc/group)
#		mount -t ${usbfs} ${usbgid:+-o devmode=0664,devgid=$usbgid}
		mount -nt ${usbfs} none /proc/bus/usb &>/dev/null
		eend $? "Failed to mount USB device filesystem"
	fi

	# Swap on loopback devices, and other weirdnesses
	ebegin "Activating (possibly) more swap"
	/sbin/swapon -a
	eend $?

}

# vim:ts=4
