# src/embedded/baselayout-lite/image/etc/init.d/modules
#!/sbin/runscript
# Copyright 1999-2005 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo/src/embedded/baselayout-lite/image/etc/init.d/Attic/modules,v 1.2 2005/03/06 06:28:08 iggy Exp $

depend() {
	need checkroot hostname
}

start() {
	# Should not fail if kernel do not have module
	# support compiled in ...
	[ -f /proc/modules ] || return 0

	# Here we should fail, as a modular kernel do need
	# depmod command ...
	if [ ! -x /sbin/depmod ]
	then
		eerror "ERROR:  system is missing /sbin/depmod !"
		return 1
	fi

	eend 0 "/etc/init.d/modules not implemented yet"




	if [ -z "${CDBOOT}" ] && touch /etc/modules.conf 2> /dev/null
	then
		ebegin "Calculating module dependencies"
		/sbin/modules-update &>/dev/null
		eend $? "Failed to calculate module dependencies"
	fi

	if [ -f /etc/modules.autoload -a ! -L /etc/modules.autoload ]
	then
		# Loop over every line in /etc/modules.autoload.
		load_modules /etc/modules.autoload
	else
		local KV="$(uname -r)"
		local KV_MAJOR="`KV_major "${KV}"`"
		local KV_MINOR="`KV_minor "${KV}"`"

		# New support for /etc/modules.autoload/kernel-$KV
		if [ "$(get_KV)" -ge "$(KV_to_int '2.5.48')" ] && \
		   [ -f /etc/modules.autoload.d/kernel-"${KV_MAJOR}.${KV_MINOR}" ]
		then
			load_modules /etc/modules.autoload.d/kernel-"${KV_MAJOR}.${KV_MINOR}"
			
		elif [ ! -f /etc/modules.autoload.d/kernel-"${KV_MAJOR}.${KV_MINOR}" ]
		then
			ewarn "Missing /etc/modules.autoload.d/kernel-${KV_MAJOR}.${KV_MINOR}"
			load_modules /etc/modules.autoload.d/kernel-2.4
		else
			load_modules /etc/modules.autoload.d/kernel-2.4
		fi
	fi

	#
	# Just in case a sysadmin prefers generic symbolic links in
	# /lib/modules/boot for boot time modules we will load these modules
	#
	if [ -n "$(modprobe -l -t boot)" ]
	then
		modprobe -a -t boot \*  &>/dev/null
	fi
}


# vim:ts=4
