# src/embedded/baselayout-lite/image/etc/init.d/net
#!/sbin/runscript
# Copyright 1999-2005 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

# much less magic here than the default Gentoo net scripts
# basically, you have to do everything yourself

depend() {
	use coldplug hotplug pcmcia usb isdn4linux wlan
}

start() {
	# start lo
	ifconfig lo 127.0.0.1 up 2>/dev/null
	/sbin/route add -net 127.0.0.0 netmask 255.0.0.0 \
		gw 127.0.0.1 dev lo 2> /dev/null

	# run whatever command the user wants to setup their network
	$NET_START_COMMAND
}

stop() {
	# take lo down
	ifconfig lo down

	# run stop command
	$NET_STOP_COMMAND
}


# vim:ts=4
