# src/embedded/baselayout-lite/image/etc/init.d/rmnologin
#!/sbin/runscript
# Copyright 1999-2005 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo/src/embedded/baselayout-lite/image/etc/init.d/Attic/rmnologin,v 1.1 2005/03/01 05:17:25 iggy Exp $

depend() {
	need localmount
}

start() {
	if [[ -f /etc/nologin.boot ]] ; then
		rm -f /etc/nologin /etc/nologin.boot &> /dev/null
	fi
}


# vim:ts=4
