# src/livecd/profiles/global/aux-files/modules
#!/sbin/runscript
# Copyright 1999-2002 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: /var/cvsroot/gentoo/src/livecd/profiles/global/aux-files/Attic/modules,v 1.1 2003/10/19 03:10:52 livewire Exp $


depend() {
	need bootmisc localmount
}

start() {
	[ -f /proc/modules ] || return
	[ -x /sbin/depmod ] || return
	ebegin "Calculating module dependencies"
	if [ -z "$CDBOOT" ]
	then
		/usr/sbin/update-modules &>/dev/null
		eend $? "Failed to calculate dependencies"
	else
		/usr/sbin/update-modules &>/dev/null
		eend
	fi
	# Loop over every line in /etc/modules.autoload.
	(cat /etc/modules.autoload; echo) | # make sure there is a LF at the end
	while read module args
	do
		case "${module}" in
			\#*|"") continue ;;
		esac
		ebegin "  Loading module ${module}"
		modprobe ${module} ${args} &>/dev/null
		eend $? "  Failed to load ${module}"
	done
	echo

	#
	# Just in case a sysadmin prefers generic symbolic links in
	# /lib/modules/boot for boot time modules we will load these modules
	#
	if [ -n "$(modprobe -l -t boot)" ]
	then
		modprobe -a -t boot \*  &>/dev/null
	fi
}


# vim:ts=4
