# users/beejay/baselayout/etc/init.d.new/consolefont
#!/sbin/runscript
# Copyright 1999-2003 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo/users/beejay/baselayout/etc/init.d.new/consolefont,v 1.1 2003/12/31 19:38:20 beejay Exp $
TEXTDOMAINDIR=/etc/gentoo/locale
TEXTDOMAIN=baselayout


depend() {
	need localmount
}

start() {
	local retval=1
	if [ -n "${CONSOLEFONT}" ]
	then
		# Getting additional parameters, ie consoletranslation
		local param=""
		local x=""
		if [ -n "${CONSOLETRANSLATION}" ]
		then
			param="-m ${CONSOLETRANSLATION}"
		fi
		
		# Actually setting font
		if [ -x /usr/bin/consolechars ]
		then
			#
			# Console-Tools section:
			#
			
			ebegin $"Setting user font"
			for x in 1 2 3 4 5 6 7 8 9 10
			do
				# Using DEVFS ?
				if [ -e /dev/.devfsd ]
				then
					/usr/bin/consolechars -f ${CONSOLEFONT} ${param} \
						--tty=/dev/vc/${x} &>/dev/null
					retval=$?
				else
					/usr/bin/consolechars -f ${CONSOLEFONT} ${param} \
						--tty=/dev/tty${x} &>/dev/null
					retval=$?
				fi
			done
	
		elif [ -x /usr/bin/setfont ]
		then
			#
			# Kbd Section
			#
			
			ebegin $"Setting user font"
			/usr/bin/setfont ${CONSOLEFONT} ${param} >/dev/null
			retval=$?
			
			for x in 1 2 3 4 5 6 7 8 9 10
			do
				# Using DEVFS ?
				if [ -e /dev/.devfsd ]
			    	then
					echo -ne "\033(K" > /dev/vc/${x} 
					retval=$?
				else
					echo -ne "\033(K" > /dev/tty${x} 
					retval=$?
				fi
			done
	
		else
			# No console program installed!
			eend 1 $"The setfont/consolechars executable not found"
			return 1
		fi
		
		eend $retval $"Failed to set user font"
	else
		ebegin $"Using the default console font"
		eend 0
	fi
	
	return $retval
}


# vim:ts=4
