# users/beejay/baselayout/etc/init.d.new/keymaps
#!/sbin/runscript
# Copyright 1999-2003 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo/users/beejay/baselayout/etc/init.d.new/keymaps,v 1.1 2003/12/31 19:38:20 beejay Exp $
TEXTDOMAINDIR=/etc/gentoo/locale
TEXTDOMAIN=baselayout


depend() {
	need localmount
}

checkconfig() {
	if [ -z "${KEYMAP}" ]
	then
		eerror $"You need to setup \$KEYMAP in /etc/rc.conf first"
		return 1
	fi
}

start() {
	local WINDOWKEYS_KEYMAP=
	
	checkconfig || return 1

	# Force linux keycodes for PPC ...
	if [ -f /proc/sys/dev/mac_hid/keyboard_sends_linux_keycodes ]
	then
		echo 1 > /proc/sys/dev/mac_hid/keyboard_sends_linux_keycodes
	fi
	
	ebegin $"Loading key mappings"
	
	# Enable unicode ...
	set -- ${KEYMAP}
	if [ "$1" = "-u" ]
	then
		/usr/bin/kbd_mode -u
	fi
	
	# New kbd support ...
	if [ -x /bin/loadkeys ]
	then
		[ "${SET_WINDOWKEYS}" = "yes" ] && WINDOWKEYS_KEYMAP="windowkeys"
		
		/bin/loadkeys -q ${WINDOWKEYS_KEYMAP} ${KEYMAP} \
		                 ${EXTENDED_KEYMAPS} >/dev/null
		
	# Old console-tools support ...
	elif [ -x /usr/bin/loadkeys ]
	then
		[ "${SET_WINDOWKEYS}" = "yes" ] && WINDOWKEYS_KEYMAP="windowkeys.inc"
		
		/usr/bin/loadkeys -q ${WINDOWKEYS_KEYMAP} ${KEYMAP} \
		                     ${EXTENDED_KEYMAPS} >/dev/null
	else
		eend 1 $"loadkeys executable not found"
		return 1
	fi
	eend $? $"Error loading key mappings"
}


# vim:ts=4
