# users/beejay/baselayout/etc/init.d.new/nscd
#!/sbin/runscript
# Copyright 1999-2003 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo/users/beejay/baselayout/etc/init.d.new/nscd,v 1.1 2003/12/31 19:38:20 beejay Exp $
TEXTDOMAINDIR=/etc/gentoo/locale
TEXTDOMAIN=baselayout


start() {
	ebegin $"Starting Name Service Cache Daemon"
	secure=`while read curline ; do
		table=${curline%:*}
		entries=${curline##$table:}
		table=${table%%[^a-z]*}
		case $table in
			passwd*|group*|hosts)
			for entry in $entries ; do
			case $entry in
				nisplus*)
					/usr/sbin/nscd_nischeck $table || \
					/echo "-S $table,yes"
				;;
			esac
			done
		;;
		esac
		done < /etc/nsswitch.conf`
	start-stop-daemon --start --quiet --exec /usr/sbin/nscd -- $secure
	eend $?
}

stop () {
	ebegin $"Shutting down Name Service Cache Daemon"
	start-stop-daemon --stop --quiet --pid /var/run/nscd.pid
	eend $?
}


# vim:ts=4
