# users/beejay/baselayout/etc/init.d.new/numlock
#!/sbin/runscript
# Copyright 1999-2003 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo/users/beejay/baselayout/etc/init.d.new/numlock,v 1.1 2003/12/31 19:38:20 beejay Exp $
TEXTDOMAINDIR=/etc/gentoo/locale
TEXTDOMAIN=baselayout


depend() {
	need localmount
}

start() {
	ebegin $"Enabling numlock on ttys"
	for tty in /dev/tty[1-8]
	do
		setleds -D +num < ${tty} &>/dev/null
	done
	eend $? $"Failed to enable numlock"
}

stop() {
	ebegin $"Disabling numlock on ttys"
	for tty in /dev/tty[1-8]
	do
		setleds -D -num < ${tty} &>/dev/null
	done
	eend $? $"Failed to disable numlock"
}


# vim:ts=4
