# src/gwn/find-cvs-adds-and-removals/find-cvs-adds-and-removals.crontab
5 0 * * Mon ./find-cvs-adds-and-removals-wrapper.sh ; ./find-cvs-adds-and-removals-mailer.sh
