# src/livecd-tools/spind.init
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo/src/livecd-tools/spind.init,v 1.2 2004/07/25 23:02:27 plasmaroo Exp $

depend() {
	before local
}

start() {
	ebegin 'Starting spin daemon'
	spind &
	sleep 1
	eend $?
}
