# src/livecd-tools/x-setup.init
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo/src/livecd-tools/x-setup.init,v 1.1 2005/02/04 13:46:11 wolf31o2 Exp $

depend() {
	after local
	before xdm
	need mkxf86config
}

start() {
	ebegin "Configuring OpenGL"
	/usr/sbin/x-setup
	eend $?
}
