# users/beejay/baselayout/etc/init.d.old/local
#!/sbin/runscript
# Copyright 1999-2003 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo/users/beejay/baselayout/etc/init.d.old/local,v 1.1 2003/12/31 20:30:00 beejay Exp $


depend() {
	after *
}

start() {
	ebegin "Starting local"

	# Add any misc programs that should be started
	# to /etc/conf.d/local.start
	if [ -e /etc/conf.d/local.start ]
	then
		source /etc/conf.d/local.start
	fi
	
	eend $? "Failed to start local."
}

stop() {
	ebegin "Stopping local"

	# Add any misc programs that should be stopped
	# to /etc/conf.d/local.stop

	if [ -e /etc/conf.d/local.stop ]
	then
		source /etc/conf.d/local.stop
	fi
	
	eend $? "Failed to stop local."
}


# vim:ts=4
