# users/beejay/baselayout/etc/init.d.old/rmnologin
#!/sbin/runscript
# Copyright 1999-2003 Gentoo Technologies, Inc.
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo/users/beejay/baselayout/etc/init.d.old/rmnologin,v 1.1 2003/12/31 20:30:00 beejay Exp $


depend() {
	need localmount
}

start() {
	if [ -f /etc/nologin.boot ]
	then
		rm -f /etc/nologin /etc/nologin.boot &>/dev/null
	fi
}


# vim:ts=4
